# drophost.service
[Service]
Type=simple
ExecStart=/usr/bin/drophost -w
Restart=always
RestartSec=10

[Install]
WantedBy=network-online.target