# drophost
#!/sbin/openrc-run

# This file is an OpenRC service for drophost.
# To install it, copy it to /etc/init.d/drophost and run: rc-update add drophost default


depend() {
  need hostname root
  use net
}

command="/root/.cargo/bin/drophost"
command_args="-w -d"
command_background="yes"
pidfile="/run/drophost.pid"